# src/tpnsfand.service
# SPDX-License-Identifier: GPL-2.0-only
[Unit]
Description=tpnsfand - ThinkPad non-standard EC fan control
After=systemd-modules-load.service
StartLimitIntervalSec=0

[Service]
Type=simple
ExecStart=/usr/local/sbin/tpnsfand
Nice=-10
OOMScoreAdjust=-1000
Restart=on-failure
RestartSec=2

CapabilityBoundingSet=CAP_SYS_NICE CAP_SYS_RESOURCE
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
PrivateNetwork=yes
PrivateMounts=yes
PrivateDevices=yes
PrivateIPC=yes
ProtectControlGroups=yes
ProtectKernelModules=yes
ProtectKernelLogs=yes
ProtectClock=yes
ProtectHostname=yes
ProtectProc=invisible
ProcSubset=pid
RestrictNamespaces=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
RestrictAddressFamilies=none
IPAddressDeny=any
LockPersonality=yes
MemoryDenyWriteExecute=yes
RemoveIPC=yes
SystemCallArchitectures=native
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM
UMask=0077
KeyringMode=private

# Uncomment and edit to override defaults (no validation is done):
#Environment=TPNSFAND_INTERVAL=
#Environment=TPNSFAND_TEMP_FAN_ON=
#Environment=TPNSFAND_TEMP_FAN_OFF=
#Environment=TPNSFAND_TEMP_CTRL_EC=
#Environment=TPNSFAND_TEMP_CTRL_HOST=
#Environment=TPNSFAND_RPM_LOW=
#Environment=TPNSFAND_RPM_HIGH=
#Environment=TPNSFAND_EC_IO=
#Environment=TPNSFAND_REG_CTRL=
#Environment=TPNSFAND_REG_SPEED=
#Environment=TPNSFAND_CTRL_HOST=
#Environment=TPNSFAND_CTRL_EC=
#Environment=TPNSFAND_SPEED_MIN=
#Environment=TPNSFAND_SPEED_MAX=
#Environment=TPNSFAND_SPEED_OFF=

[Install]
WantedBy=multi-user.target
